# sv/st2chatops/run
#!/bin/bash

sv start prestart || sleep 1 || exit 1

source /etc/envvars

export ST2_AUTH_TOKEN=`st2 auth -t -p $ST2_AUTH_PASSWORD $ST2_AUTH_USERNAME`
export ST2_API_KEY=`st2 apikey create -k`

exec 2>&1
exec /opt/stackstorm/chatops/bin/hubot
