# sv/st2notifier/run
#!/bin/bash

sv start prestart || sleep 1 || exit 1

source /etc/envvars

exec 2>&1
exec  /opt/stackstorm/st2/bin/st2notifier --config-dir /etc/st2/
